# homemind-pro/rootfs/etc/s6-overlay/s6-rc.d/install-integration/run
#!/command/with-contenv bash
# Auto-install the Home Mind HA integration (custom component)
# so users don't need HACS or manual file copying.

SOURCE="/opt/home_mind_integration"
TARGET="/config/custom_components/home_mind"

if [ ! -d /config ]; then
    echo "[integration] /config not mounted — skipping integration install"
    exit 0
fi

# Check if files need updating by comparing manifest versions
NEEDS_UPDATE=false
if [ ! -d "$TARGET" ]; then
    NEEDS_UPDATE=true
    echo "[integration] Installing Home Mind integration..."
elif [ -f "$SOURCE/manifest.json" ] && [ -f "$TARGET/manifest.json" ]; then
    SOURCE_VER=$(jq -r '.version' "$SOURCE/manifest.json" 2>/dev/null || echo "")
    TARGET_VER=$(jq -r '.version' "$TARGET/manifest.json" 2>/dev/null || echo "")
    if [ "$SOURCE_VER" != "$TARGET_VER" ]; then
        NEEDS_UPDATE=true
        echo "[integration] Updating Home Mind integration (${TARGET_VER} → ${SOURCE_VER})..."
    fi
fi

if [ "$NEEDS_UPDATE" = "true" ]; then
    mkdir -p "$TARGET"
    cp -f "$SOURCE"/*.py "$TARGET/"
    cp -f "$SOURCE"/*.json "$TARGET/"
    echo "[integration] Home Mind integration installed at $TARGET"
    echo "[integration] HA restart may be required to pick up the integration"
else
    echo "[integration] Home Mind integration is up to date"
fi
